# services/kakoune-server/run
#!/usr/bin/env bash

cd $HOME

export XDG_RUNTIME_DIR=/run/user/$(id -u)

if [ -n "$(kak -l)" ]; then
    kak -clear
fi

exec zsh -lc "kak -s $USER -d"
