# services/emacs-server/run
#!/bin/sh

exec zsh -i -c "emacs --fg-daemon"
