# services/emacs-server/run
#!/bin/sh

cd $HOME
exec zsh -i -c "emacs --fg-daemon"
